# synapse-node/systemd/synapse-node.service
[Unit]
Description=Synapse Network Node
Documentation=https://docs.synapse.network
After=network.target nvidia-persistenced.service docker.service
Wants=network.target
Requires=nvidia-persistenced.service

[Service]
Type=simple
User=synapse
Group=synapse

# Working directory
WorkingDirectory=/opt/synapse

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="SYNAPSE_LOG_LEVEL=INFO"
Environment="SYNAPSE_MODEL_CACHE_DIR=/var/lib/synapse/models"
Environment="SYNAPSE_PROOF_CIRCUIT_PATH=/var/lib/synapse/circuits"

# Load environment from file
EnvironmentFile=-/etc/synapse/node.env

# Main executable
ExecStart=/opt/synapse/venv/bin/synapse-node start

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
# Adjust based on your hardware
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/synapse /var/log/synapse
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# GPU access
# Note: NVIDIA devices need to be accessible
DeviceAllow=/dev/nvidia* rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-uvm rw
DeviceAllow=/dev/nvidia-uvm-tools rw

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=synapse-node

[Install]
WantedBy=multi-user.target
